# scripts/files/camera-streamer/camera-streamer-usb.path
[Unit]
Description=camera-streamer usb autolaunch

[Path]
PathExistsGlob=/dev/v4l/by-id/usb-*-video-index0
Unit=camera-streamer-usb.service

[Install]
WantedBy=multi-user.target
