# scripts/files/camera-streamer/camera-streamer-usb@.path
[Unit]
Description=camera-streamer %i autolaunch
After=network.target

[Path]
PathExistsGlob=/dev/v4l/by-id/usb-*-video-index0

[Install]
WantedBy=multi-user.target
